# stroke-memorize.service
[Unit]
Description=stroke-memorize kanji stroke order app
After=network.target

[Service]
Type=simple
User=stroke-memorize
WorkingDirectory=/var/lib/stroke-memorize/stroke-order-memorize
EnvironmentFile=/var/lib/stroke-memorize/stroke-order-memorize/.env
ExecStart=/usr/bin/uv run uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/stroke-memorize/stroke-order-memorize %h/.cache/uv
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
